# .claude/skills/claap-weekly-recap/scripts/launchd.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Claap Weekly Recap — macOS launchd schedule

  Install:
    1. Replace every <PLACEHOLDER> below with your absolute paths and username.
    2. Copy this file to ~/Library/LaunchAgents/com.yalc.claap-weekly-recap.plist
    3. launchctl load ~/Library/LaunchAgents/com.yalc.claap-weekly-recap.plist
    4. launchctl start com.yalc.claap-weekly-recap  (one-off manual fire to test)

  Schedule: every Sunday at 08:00 local time (Weekday=0 in launchd convention
  where Sunday=0, Saturday=6). Change the Hour/Minute/Weekday values to taste.

  Uninstall: launchctl unload ~/Library/LaunchAgents/com.yalc.claap-weekly-recap.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yalc.claap-weekly-recap</string>

    <key>ProgramArguments</key>
    <array>
        <string><ABSOLUTE_PATH_TO_run.sh></string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/claap-weekly-recap.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claap-weekly-recap.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/<YOUR_MACOS_USERNAME></string>
    </dict>
</dict>
</plist>
